# Scripts/cryptoalgo.crontab
# Run CryptoAlgo daily at 10am
0 10 * * * cd /home/ec2-user/CryptoAlgo && /usr/bin/python3 CryptoPriceDBGateway.py > /dev/null
